# agent/cap-sync.service
[Unit]
Description=Cap Tracker Agent - Sync development activity data
After=network.target

[Service]
Type=oneshot
ExecStart=__NPX__ tsx __PROJECT_DIR__/agent/src/cli.ts sync
WorkingDirectory=__PROJECT_DIR__
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
