# agent/cap-sync.timer
[Unit]
Description=Run Cap Tracker sync during business hours + late evening

[Timer]
OnCalendar=Mon..Fri 08,10,12,14,16,18,23:00
OnCalendar=Sat,Sun 12,23:00
Persistent=true

[Install]
WantedBy=timers.target
