# agent/cap-generate.service
[Unit]
Description=Cap Tracker - Generate AI daily entries for yesterday
After=network.target cap-sync.service

[Service]
Type=oneshot
# First sync to get latest data, then generate entries for yesterday
ExecStartPre=__NPX__ tsx __PROJECT_DIR__/agent/src/cli.ts sync
ExecStart=__NPX__ tsx __PROJECT_DIR__/agent/src/cli.ts generate
WorkingDirectory=__PROJECT_DIR__
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
